# config/my-app.path
[Unit]
Description=start my app, wow!

[Install]
WantedBy=multi-user.target

[Path]
PathExists=/srv/my-app/current
